# etc/init.d/ananicy-cpp
#!/sbin/openrc-run

description="Ananicy-Cpp - ANother Auto NICe daemon in C++"

command="/usr/bin/ananicy-cpp"
supervisor=supervise-daemon
command_args_foreground="start"

respawn_delay=10
respawn_max=0
supervise_daemon_args="--nicelevel -5 --oom-score-adj -999"

extra_started_commands="reload"

depend() {
    need localmount
}

reload() {
    ebegin "Reloading ananicy-cpp rules"
    $command --reload
    eend $?
}
